# container/supervisord/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/www/html/log/supervisor/supervisord.log
childlogdir=/var/www/html/log/supervisor

[program:rsyslog]
command=/usr/sbin/rsyslogd -n
numprocs=1
autostart=true
autorestart=true
priority=1

[program:mysql]
command=mysql-init.sh
autostart=true
autorestart=true
priority=10

[program:php-fpm]
command=/usr/sbin/php-fpm8.2 -c /etc/php/8.2/fpm
autostart=true
autorestart=true
exitcodes=0,2,70
stopsignal=QUIT
priority=10

[program:init]
command=bash -c 'sleep 10 && /init.sh'
startretries=0
autostart=true
autorestart=false
